# briefing-bot/systemd/briefing-bot.timer
[Unit]
Description=Run AI briefing bot every day at 08:50 CST

[Timer]
OnCalendar=*-*-* 08:50:00 Asia/Shanghai
Persistent=true

[Install]
WantedBy=timers.target
