# briefing-bot/systemd/briefing-bot-reply.service
[Unit]
Description=Feishu app bot reply service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/wuchao/ai/my-vault
ExecStart=/home/wuchao/ai/my-vault/briefing-bot/.venv/bin/python /home/wuchao/ai/my-vault/briefing-bot/scripts/feishu_reply.py
Restart=always
RestartSec=5
User=wuchao
Environment=HOME=/home/wuchao
Environment=PATH=/home/wuchao/ai/my-vault/briefing-bot/.venv/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
